# gpa-updater.service.example
[Unit]
Description=gpa-updater

[Service]
Type=oneshot
ExecStart=/home/pi/gpa-updater.py